# ansible/roles/run_zeilvoortgang/templates/zeilvoortgang.service.j2
[Unit]
Description=Java zeilvoortgang Service
After=network.target

[Service]
User=root
WorkingDirectory=/opt/zeilvoortgang
ExecStart=/usr/bin/java -jar /opt/zeilvoortgang/zeilvoortgang.jar
Restart=always
RestartSec=10
Environment=SPRING_PROFILES_ACTIVE=prod

[Install]
WantedBy=multi-user.target
